# files/etc/cron.d/dnscrypt
@reboot root /etc/dnscrypt-wrapper/start.sh
0 * * * * root /etc/dnscrypt-wrapper/start.sh