# files/etc/cron.d/recursor
0 3 * * * root dig @178.63.145.230 . NS > /etc/powerdns/opennic.hint.tmp && fgrep ";; ->>HEADER<<- opcode: QUERY, status: NOERROR, id: " /etc/powerdns/opennic.hint.tmp && mv /etc/powerdns/opennic.hint.tmp /etc/powerdns/opennic.hint && rec_control reload-zones > /dev/null 2>&1
* * * * * root /etc/init.d/pdns-recursor start
@reboot root /etc/init.d/pdns-recursor start